# .config/systemd/user/pomodoro.service
[Unit]
Description=Pomodoro timer

[Service]
Type=simple
ExecStart=%h/.local/bin/pomodoro

[Install]
WantedBy=default.target
